# board/common/rootfs-overlay/etc/systemd/system/iwd.service.d/install-keys.conf
[Service]
StateDirectory=iwd
StateDirectoryMode=0700
ExecStartPre=/bin/find /boot/ -maxdepth 1 -name '*.psk' -type f -exec /bin/cp -p -f {} /var/lib/iwd/ +
