# board/common/rootfs-overlay/etc/systemd/system/bluetooth.service.d/compat.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --compat
