# board/common/rootfs-overlay/etc/systemd/system/bluetooth.service.d/discover.conf
[Service]
ExecStartPost=-/usr/bin/sdptool add SP
ExecStartPost=/usr/bin/hciconfig hci0 piscan
