# board/common/rootfs-overlay/etc/systemd/system/pulseaudio.service.d/system.conf
[Service]
ExecStart=
ExecStart=/usr/bin/pulseaudio --system --daemonize=no \
	--disable-shm --disallow-exit --exit-idle-time=-1

